# deploy/polymarket-bot.service
[Unit]
Description=PolySignal — Polymarket Monitoring Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-bot
ExecStart=/opt/polymarket-bot/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Give it time to shut down cleanly so "Bot stopped." message is sent
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
